# files/profile.d/aliases.bash
# Add colors!
alias grep="grep --color=auto"
alias less="less --RAW-CONTROL-CHARS"
alias ls="ls --color=auto"

alias l="ls -la"
alias ..="cd .."
alias ...="cd ../.."
