# files/profile.d/history.bash
# Ignore duplicate lines
export HISTCONTROL=ignoredups

# append to the history file, don't overwrite it
shopt -s histappend

# Number of commands to store in history
export HISTSIZE=65536
export HISTFILESIZE=65536
